# systemd/hcn-init.service.in
[Unit]
Description=hybrid virtual network scan and config
After=NetworkManager.service
Requisite=NetworkManager.service
PartOf=NetworkManager.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@sbindir@/hcnmgr -s

[Install]
WantedBy=NetworkManager.service
