# systemd/hcn-init.service.in
[Unit]
Description=hybrid virtual network scan and config
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=@sbindir@/hcnmgr -s

[Install]
WantedBy=multi-user.target
